# cron/crontab
0 * * * * /usr/local/cam/bin/build_time_lapse.pl --trickle --silent
*/15 * * * * /usr/local/cam/bin/gap_check.pl
